# pkg/chumby-player/DEBIAN/prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ]; then
    systemctl disable --now chumby-player.service || true
fi
