# pkg/chumby-player/chumby-player.service
# Chumby kiosk: cage on tty1/seat0 as user pi, output on the SPI TFT.
#
# PAMName=login gives the service a real logind session (seat0 +
# XDG_RUNTIME_DIR), which is what lets cage open the DRM device as a
# plain user and lets audio reach pi's PipeWire without any bridging.
#
# The display is not named here: the launcher's --kiosk mode finds the
# SPI TFT by its stable by-path name on any Pi model and exports
# WLR_DRM_DEVICES before exec'ing cage (cardN numbering is NOT stable —
# the TFT was card1 on 2026-07-06, card0 after the next reboot).
# Override point: /etc/default/chumby-player.

[Unit]
Description=Chumby control panel kiosk (cage + ruffle)
After=systemd-user-sessions.service getty@tty1.service
Conflicts=getty@tty1.service

[Service]
Type=simple
User=pi
PAMName=login
TTYPath=/dev/tty1
StandardInput=tty-fail
StandardOutput=journal
StandardError=journal
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
UtmpIdentifier=tty1
UtmpMode=user
StateDirectory=chumby
Environment=WLR_BACKENDS=drm,libinput
Environment=WLR_RENDERER=pixman
EnvironmentFile=-/etc/default/chumby-player
# Fails with the copy-these-files instructions in `systemctl status`
# even when cage cannot start at all (no DRM device, wrong overlay).
ExecStartPre=/usr/bin/chumby-player-run --check
ExecStart=/usr/bin/chumby-player-run --kiosk
Restart=on-failure
RestartSec=3

[Install]
# multi-user, not graphical: the Lite image (the reference OS) never
# reaches graphical.target, so the kiosk would never start at boot —
# found on the 3rd vanilla-card install, 2026-07-14. Desktop-variant
# images conflict with their display manager over seat0 either way.
WantedBy=multi-user.target
