# pkg/chumby-player/chumby-widget-channel.service
# Rebuild the widget channel (/xml/profiles) from the installed widgets
# before the panel starts, so an added or removed widget shows up. Runs as
# pi with StateDirectory=chumby, matching chumby-player.service, so it
# writes the same pi-owned /var/lib/chumby/fixtures tree the panel reads.
#
# Pulled in by chumby-player.service (Wants=), not enabled on its own. A
# missing widgets tree is seeded here first, mirroring chumby-player-run's
# first-run seed, since either service may be the first to touch the tree.

[Unit]
Description=Generate the chumby widget channel from installed widgets
After=local-fs.target
Before=chumby-player.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=pi
StateDirectory=chumby
EnvironmentFile=-/etc/default/chumby-player
Environment=CHUMBY_STATE=/var/lib/chumby
ExecStartPre=/bin/sh -c 'test -d "$CHUMBY_STATE/fixtures" || { cp -a /usr/share/chumby-player/fixtures "$CHUMBY_STATE/fixtures.seed" && mv "$CHUMBY_STATE/fixtures.seed" "$CHUMBY_STATE/fixtures"; }'
ExecStart=/bin/sh -c 'exec /usr/bin/chumby-widget-channel --fixtures "$CHUMBY_STATE/fixtures" --quiet'

[Install]
WantedBy=graphical.target
